# deployment/cron/vm_project_watchdog.cron
# VM Project v7.20 Watchdog Cron
# 초등학생 설명: 1분마다 서버가 살아 있는지 확인해요.
* * * * * /home/$USER/vm_project/scripts/watchdog.sh
